# home/.crontab
*/5 5-21 * * *	PATH=~/bin:$PATH gitmap fetch >> ~/cronlog 2>&1
0 22 * * *	PATH=~/bin:$PATH gitmap echo gc >> ~/cronlog 2>&1
